# files/root/var/spool/cron/root-crons.cron
#Mins Hours   Days  Months  Day of the week
@reboot                                       /usr/bin/updatedb
@reboot                                       find /root/.cache -mtime +7 -exec rm -f {} \;
